# cron.daily/etc_git
#!/bin/sh

cd /etc
git add .
git commit -sa -m "update `date`"
